# .config/autostart/redshift.desktop
[Desktop Entry]
Exec=redshift
Icon=system-run
Path=
Terminal=false
Type=Application
